# ansible/roles/monitoring/templates/health_check.service.j2
[Unit]
Description=wireguard health check

[Service]
Type=simple
User=root
ExecStart=/root/health_check.py --interface {{ wireguard.network_interface }}

[Install]
WantedBy=multi-user.target
