# roles/monitoring/files/health_check.service
[Unit]
Description=wireguard health check

[Service]
Type=simple
User=root
ExecStart=/root/health_check.py

[Install]
WantedBy=multi-user.target
